# data/etc/supervisord.conf
[supervisord]
user        = root
nodaemon    = true
logfile     = /var/log/supervisord.log
pidfile     = /var/run/supervisord.pid
childlogdir = /var/log/

[program:httpd]
command                 = /usr/bin/httpd.py
autostart               = true
autorestart             = true
stdout_logfile          = /dev/stdout
stdout_logfile_maxbytes = 0
stdout_events_enabled   = true
stderr_logfile          = /dev/stderr
stderr_logfile_maxbytes = 0
stderr_events_enabled   = true

[program:update-metrics]
command                 = /usr/bin/update-metrics.sh
autostart               = true
autorestart             = true
stdout_logfile          = /dev/stdout
stdout_logfile_maxbytes = 0
stdout_events_enabled   = true
stderr_logfile          = /dev/stderr
stderr_logfile_maxbytes = 0
stderr_events_enabled   = true
